# crontab
*/5 * * * * /data/cron_task.sh /data/chart.py
